# profile.d/bash-env.sh
# bash-env.sh

export BASH_ENV=$HOME/.bashrc

