# profile.d/java.sh
# java.sh

export JAVA_HOME=`/usr/libexec/java_home -version 1.7`

export MAVEN_OPTS=-Xmx1024m
export MAVEN_HOME=~/.m2
