# profile.d/scala.sh
# scala.sh

export SCALA_HOME=/usr/local/Cellar/scala/2.11.7

export PATH=$SCALA_HOME/bin:$PATH
